# debian/ukui-group-manager/DEBIAN/postinst
#!/bin/bash

# 启动服务
systemctl enable ukui-group-manager.service
systemctl start  ukui-group-manager.service

